# ocr_worker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:llama-server]
command=/app/llama_cpp_bin/llama-server --host 127.0.0.1 --port 8090 -m %(ENV_OCR_VL_MODEL_PATH)s --mmproj %(ENV_OCR_VL_MMPROJ_PATH)s -c %(ENV_OCR_VL_CTX_SIZE)s -cb -b %(ENV_OCR_VL_BATCH_SIZE)s -ub %(ENV_OCR_VL_UBATCH_SIZE)s -ngl %(ENV_N_GPU_LAYERS)s -t %(ENV_OCR_VL_THREADS)s -fa on --parallel %(ENV_OCR_VL_PARALLEL)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ocr-api]
command=python -m uvicorn ocr_worker.api:app --host 0.0.0.0 --port 8091
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
